# livereduce.service
[Unit]
Description=Live processing service
StartLimitInterval=8640
StartLimitBurst=1000

[Service]
WorkingDirectory=/tmp
User=snsdata
ExecStart=/usr/bin/livereduce.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
