# livereduce.service
[Unit]
Description=Live processing service

[Service]
WorkingDirectory=/tmp
User=snsdata
ExecStart=/usr/bin/python /usr/bin/livereduce.py
Restart=always

[Install]
WantedBy=multi-user.target
